# scripts/com.dev-assistant.burnout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dev-assistant.burnout</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/dev-assistant/scripts/burnout-radar.sh</string>
    </array>

    <!-- 매주 월요일 09:00 실행 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/dev-assistant/logs/launchd-burnout.log</string>

    <key>StandardErrorPath</key>
    <string>/path/to/dev-assistant/logs/launchd-burnout.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/your-username</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
